# util/AutoStopper/autostopper.service
[Service]
Type=simple

ExecStart=/home/stemiaa/AutoStopper/AutoStopper.py

[Install]
WantedBy=multi-user.target
